# alsaconnect.service
[Unit]
Description=ALSA MIDI Connector

[Service]
ExecStart=/usr/local/midirouter/connect-puredata.sh

[Install]
WantedBy=multi-user.target
