# midirouter.service
[Unit]
Description=MIDI router

[Service]
ExecStart=/usr/local/midirouter/puredata.sh

[Install]
WantedBy=multi-user.target
